# termisol-auto-commit.service
[Unit]
Description=Termisol Auto-Commit Service
After=network.target

[Service]
Type=oneshot
ExecStart=/path/to/termisol/scripts/auto_commit.sh
User=your_username
WorkingDirectory=/path/to/termisol
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
