# scripts/termisol-user.service
[Unit]
Description=Termisol Background Service
Documentation=https://github.com/termisol/termisol
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/termisol --background
Restart=always
RestartSec=5
Environment=DISPLAY=:0
Environment=XDG_SESSION_TYPE=x11

# Security restrictions
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.termisol %h/.local/share/termisol
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
DevicePolicy=closed

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=default.target